# scripts/launchd/com.macrotracking.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macrotracking.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tranquangminhvu/Vĩ mô Mỹ Tracking/scripts/run_daily_macro.sh</string>
    </array>

    <!-- Thứ 2..6 (Weekday 1=Mon..5=Fri), 14:00 giờ máy (Vancouver) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/tranquangminhvu/Vĩ mô Mỹ Tracking/data/automation_logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tranquangminhvu/Vĩ mô Mỹ Tracking/data/automation_logs/launchd.err.log</string>
</dict>
</plist>